# root/etc/systemd/system/mnt-media.automount
[Unit]
Description=TrueNAS media share mount

[Automount]
Where=/mnt/media

[Install]
WantedBy=multi-user.target
